# chef/cookbooks/application/templates/default/api_server.app_daemon.erb
#!/bin/sh
# Comments to support chkconfig on Red Hat Linux
# chkconfig: 2345 64 36
# Description: <%= node.default['application']['api_server']['app']['name'] %>

start () {
    echo "Starting <%= node.default['application']['api_server']['app']['name'] %>... "
    daemon --name='<%= node.default['application']['api_server']['app']['name'] %>' \
           --command='<%= node.default['application']['virtual_env'] %>/bin/python <%= node.default['application']['sources_dir'] %>/interface_server/manage.py runfcgi daemonize=false method=threaded host=<%= node.default['application']['api_server']['app']['fastcgi_host'] %> port=<%= node.default['application']['api_server']['app']['fastcgi_port'] %> workdir=<%= node.default['application']['work_dir'] %> outlog=<%= node.default['application']['api_server']['app']['output_log'] %> errlog=<%= node.default['application']['api_server']['app']['error_log'] %> pidfile=<%= node.default['application']['api_server']['app']['pid_file'] %>' \
           --env='PYTHONPATH=<%= node.default['application']['python_env']['PYTHONPATH'] %>' \
           --env='DJANGO_SETTINGS_MODULE=<%= node.default['application']['python_env']['DJANGO_SETTINGS_MODULE'] %>' \
           --env='PYTHON_EGG_CACHE=<%= node.default['application']['python_env']['PYTHON_EGG_CACHE'] %>' \
           --pidfile='<%= node.default['application']['api_server']['app']['daemon']['pid_file'] %>' \
           --user='<%= node.default['application']['api_server']['user'] %>' \
           --chdir='<%= node.default['application']['work_dir'] %>' \
           --unsafe \
           --respawn \
           --errlog='<%= node.default['application']['api_server']['app']['daemon']['error_log'] %>' \
           --dbglog='<%= node.default['application']['api_server']['app']['daemon']['debug_log'] %>' \
           --verbose
    return $?
}

stop () {
    echo "Stopping <%= node.default['application']['api_server']['app']['name'] %>... "
    daemon --name='<%= node.default['application']['api_server']['app']['name'] %>' \
           --pidfile='<%= node.default['application']['api_server']['app']['daemon']['pid_file'] %>' \
           --stop \
           --verbose
    return $?
}


restart () {
    echo "Restarting <%= node.default['application']['api_server']['app']['name'] %>..."
    daemon --name='<%= node.default['application']['api_server']['app']['name'] %>' \
           --pidfile='<%= node.default['application']['api_server']['app']['daemon']['pid_file'] %>' \
           --restart \
           --verbose
    return $?
}


status () {
    daemon --name='<%= node.default['application']['api_server']['app']['name'] %>' \
           --pidfile='<%= node.default['application']['api_server']['app']['daemon']['pid_file'] %>' \
           --running \
           --verbose
    return $?
}


case "$1" in
    start)
        start
    ;;
    status)
        status
    ;;
    stop)
        stop
    ;;
    restart)
        restart
    ;;
    *)
        echo $"Usage: <%= node.default['application']['api_server']['app']['name'] %> {start|status|stop|restart}"
        exit 3
    ;;
esac

exit $RETVAL
