# chef/cookbooks/application/templates/default/api_serving.app_daemon.erb
#!/bin/sh
# Comments to support chkconfig on Red Hat Linux
# chkconfig: 2345 64 36
# Description: <%= node.default['application']['api_serving']['name'] %>

start () {
    echo "Starting <%= node.default['application']['api_serving']['app']['name'] %>... "

    daemon --name='<%= node.default['application']['api_serving']['app']['name'] %>' \
           --command='<%= node.default['application']['virtual_env'] %>/bin/python <%= node.default['application']['sources_dir'] %>/interface_server/manage.py runfcgi method=threaded host=<%= node.default['application']['api_serving']['app']['fastcgi_host'] %> port=<%= node.default['application']['api_serving']['app']['fastcgi_port'] %> workdir=<%= node.default['application']['work_dir'] %> outlog=<%= node.default['application']['api_serving']['app']['output_log'] %> errlog=<%= node.default['application']['api_serving']['app']['error_log'] %> pidfile=<%= node.default['application']['api_serving']['app']['pid_file'] %>' \
	   --env="PYTHONPATH=<%= node.default['application']['sources_dir'] %>:<%= node.default['application']['sources_dir'] %>/interface_server:<%= node.default['application']['sources_dir'] %>/gen/py" \
	   --env="DJANGO_SETTINGS_MODULE=interface_server.settings" \
           --pidfile='<%= node.default['application']['api_serving']['app']['daemon']['pid_file'] %>' \
           --user='<%= node.default['application']['api_serving']['user'] %>' \
           --chdir='<%= node.default['application']['work_dir'] %>' \
           --unsafe \
           --respawn \
           --errlog='<%= node.default['application']['api_serving']['app']['daemon']['error_log'] %>' \
           --dbglog='<%= node.default['application']['api_serving']['app']['daemon']['debug_log'] %>' \
           --verbose
    return $?
}

stop () {
    echo "Stopping <%= node.default['application']['api_serving']['app']['name'] %>... "

    daemon --name='<%= node.default['application']['api_serving']['app']['name'] %>' \
           --pidfile='<%= node.default['application']['api_serving']['app']['daemon']['pid_file'] %>' \
           --stop \
           --verbose
    return $?
}


restart () {
    echo "Restarting <%= node.default['application']['api_serving']['app']['name'] %>..."

    daemon --name='<%= node.default['application']['api_serving']['app']['name'] %>' \
           --pidfile='<%= node.default['application']['api_serving']['app']['daemon']['pid_file'] %>' \
           --restart  \
           --verbose
    return $?
}


status () {
    daemon --name='<%= node.default['application']['api_serving']['app']['name'] %>' \
           --pidfile='<%= node.default['application']['api_serving']['app']['daemon']['pid_file'] %>' \
           --running \
           --verbose
    return $?
}


case "$1" in
    start)
        start
    ;;
    status)
        status
    ;;
    stop)
        stop
    ;;
    restart)
        restart
    ;;
    *)
        echo $"Usage: <%= node.default['application']['api_serving']['app']['name'] %> {start|status|stop|restart}"
        exit 3
    ;;
esac

exit $RETVAL
