# chef/cookbooks/application/templates/default/explorer.explorer_daemon.erb
#!/bin/sh
# Comments to support chkconfig on Red Hat Linux
# chkconfig: 2345 64 36
# Description: <%= node.default['application']['explorer']['explorer']['name'] %>

start () {
    echo "Starting <%= node.default['application']['explorer']['explorer']['name'] %>... "
    daemon --name='<%= node.default['application']['explorer']['explorer']['name'] %>' \
           --command='<%= node.default['application']['virtual_env'] %>/bin/python <%= node.default['application']['sources_dir'] %>/explorer/server.py --sleep_sec=<%= node.default['application']['explorer']['explorer']['sleep_sec'] %> --fetcher_host=<%= node.default['application']['explorer']['fetcher']['host'] %> --fetcher_port=<%= node.default['application']['explorer']['fetcher']['port'] %> --photo_save_host=<%= node.default['application']['explorer']['photo_save']['host'] %> --photo_save_port=<%= node.default['application']['explorer']['photo_save']['port'] %> --counters_log_path=<%= node.default['application']['explorer']['explorer']['counters_log'] %> --log_path=<%= node.default['application']['explorer']['explorer']['log'] %> --pidfile_path=<%= node.default['application']['explorer']['explorer']['pid_file'] %>' \
           --env='PYTHONPATH=<%= node.default['application']['python_env']['PYTHONPATH'] %>' \
           --env='DJANGO_SETTINGS_MODULE=<%= node.default['application']['python_env']['DJANGO_SETTINGS_MODULE'] %>' \
           --env='PYTHON_EGG_CACHE=<%= node.default['application']['python_env']['PYTHON_EGG_CACHE'] %>' \
           --pidfile='<%= node.default['application']['explorer']['explorer']['daemon']['pid_file'] %>' \
           --user='<%= node.default['application']['explorer']['user'] %>' \
           --chdir='<%= node.default['application']['work_dir'] %>' \
           --unsafe \
           --respawn \
           --errlog='<%= node.default['application']['explorer']['explorer']['daemon']['error_log'] %>' \
           --dbglog='<%= node.default['application']['explorer']['explorer']['daemon']['debug_log'] %>' \
           --verbose
    return $?
}

stop () {
    echo "Stopping <%= node.default['application']['explorer']['explorer']['name'] %>... "
    daemon --name='<%= node.default['application']['explorer']['explorer']['name'] %>' \
           --pidfile='<%= node.default['application']['explorer']['explorer']['daemon']['pid_file'] %>' \
           --stop \
           --verbose
    return $?
}


restart () {
    echo "Restarting <%= node.default['application']['explorer']['explorer']['name'] %>..."
    daemon --name='<%= node.default['application']['explorer']['explorer']['name'] %>' \
           --pidfile='<%= node.default['application']['explorer']['explorer']['daemon']['pid_file'] %>' \
           --restart \
           --verbos
    return $?
}


status () {
    daemon --name='<%= node.default['application']['explorer']['explorer']['name'] %>' \
           --pidfile='<%= node.default['application']['explorer']['explorer']['daemon']['pid_file'] %>' \
           --running \
           --verbose
    return $?
}


case "$1" in
    start)
        start
    ;;
    status)
        status
    ;;
    stop)
        stop
    ;;
    restart)
        restart
    ;;
    *)
        echo $"Usage: <%= node.default['application']['explorer']['explorer']['name'] %> {start|status|stop|restart}"
        exit 3
    ;;
esac

exit $RETVAL
