# chef/cookbooks/application/templates/default/redirect_server.frontend_daemon.erb
#!/bin/sh
# Comments to support chkconfig on Red Hat Linux
# chkconfig: 2345 64 36
# Description: <%= node.default['application']['redirect_server']['frontend']['name'] %>

start () {
    echo "Starting <%= node.default['application']['redirect_server']['frontend']['name'] %>... "
    daemon --name='<%= node.default['application']['redirect_server']['frontend']['name'] %>' \
           --command='lighttpd -D -f <%= node.default['application']['redirect_server']['frontend']['config'] %>' \
           --pidfile='<%= node.default['application']['redirect_server']['frontend']['daemon']['pid_file'] %>' \
           --user='<%= node.default['application']['redirect_server']['user'] %>' \
           --chdir='<%= node.default['application']['work_dir'] %>' \
           --safe \
           --respawn \
           --errlog='<%= node.default['application']['redirect_server']['frontend']['daemon']['error_log'] %>' \
           --dbglog='<%= node.default['application']['redirect_server']['frontend']['daemon']['debug_log'] %>' \
           --verbose
    return $?
}

stop () {
    echo "Stopping <%= node.default['application']['redirect_server']['frontend']['name'] %>... "
    daemon --name='<%= node.default['application']['redirect_server']['frontend']['name'] %>' \
           --pidfile='<%= node.default['application']['redirect_server']['frontend']['daemon']['pid_file'] %>' \
           --stop \
           --verbose
    return $?
}


restart () {
    echo "Restarting <%= node.default['application']['redirect_server']['frontend']['name'] %>..."
    daemon --name='<%= node.default['application']['redirect_server']['frontend']['name'] %>' \
           --pidfile='<%= node.default['application']['redirect_server']['frontend']['daemon']['pid_file'] %>' \
           --restart  \
           --verbose
    return $?
}


status () {
    daemon --name='<%= node.default['application']['redirect_server']['frontend']['name'] %>' \
           --pidfile='<%= node.default['application']['redirect_server']['frontend']['daemon']['pid_file'] %>' \
           --running \
           --verbose
    return $?
}


case "$1" in
    start)
        start
    ;;
    status)
        status
    ;;
    stop)
        stop
    ;;
    restart)
        restart
    ;;
    *)
        echo $"Usage: <%= node.default['application']['redirect_server']['frontend']['name'] %> {start|status|stop|restart}"
        exit 3
    ;;
esac

exit $RETVAL
